# roles/restic/templates/restic-backup.timer.j2
[Unit]
Description=Timer for Restic Backup

[Timer]
OnCalendar={{ restic.schedule.backup }}
Persistent=true

[Install]
WantedBy=timers.target
