# roles/restic/templates/restic-check.service.j2
[Unit]
Description=Restic Repository Integrity Check
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
WorkingDirectory={{ restic.path }}
ExecStart=/usr/bin/docker compose run --rm --entrypoint "restic check --read-data-subset={{ restic.check.subset | replace('%', '%%') }}" worker

[Install]
WantedBy=multi-user.target
