# roles/restic/templates/restic-check.timer.j2
[Unit]
Description=Weekly Schedule for Restic Check

[Timer]
OnCalendar={{ restic.schedule.check }}
Persistent=true

[Install]
WantedBy=timers.target
